# FokinCPL/fokincpl.service
[Unit]
Description=My fokin CPL are bugged

[Service]
Type=simple
ExecStart=/bin/bash "PATH OF THE SCRIPT"

[Install]
WantedBy=multi-user.target
